# root/etc/services.d/bot/run
#!/usr/bin/with-contenv bash

cd /app/gmail || exit
exec \
	s6-setuidgid abc node index.js
